# scripts/container/buildah/assets/s6-rc.d/orbuculum-web-log-pre/up
if { mkdir -p /var/log/orbuculum-web }
if { chown nobody:nogroup /var/log/orbuculum-web }
chmod 02755 /var/log/orbuculum-web
